# supervisor/app.conf

[program:camera]
environment=DISPLAY=":1"
command=python3 camera_app.py
